# dist/initscripts/rc.local
#
# Add this line to /etc/rc.local
#

/home/upnp/sys/boot.sh &

